# cronjobs/cron.hourly/yay_pkglist
#!/bin/sh

pacman -Qqe > /home/gchamon/Documents/pkglist.txt
chown gchamon:gchamon /home/gchamon/Documents/pkglist.txt

