# cronjobs/cron.hourly/yay_pkglist
#!/bin/sh

pacman -Qqe >/etc/pkglist.txt
